# files/etc/supervisor/conf.d/deluge-web.conf
[program:deluge-web]
command = /usr/bin/deluge-web -c /data
