# files/etc/supervisor/conf.d/deluged.conf
[program:deluged]
command = /usr/bin/deluged -c /data -d
